# cron/health_check.cron
*/5 * * * * /home/cao/linux-web-monitoring/scripts/health_check.sh
